# supervisor/playwright-mcp.conf
[program:playwright-mcp-gateway]
command=/usr/local/bin/start-playwright-mcp.sh
directory=/tmp
autostart=true
autorestart=true
priority=300
stderr_logfile=/var/log/supervisor/playwright-mcp-gateway.err.log
stdout_logfile=/var/log/supervisor/playwright-mcp-gateway.out.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
user=root
startsecs=10
startretries=5
killasgroup=true
stopasgroup=true
environment=PLAYWRIGHT_MCP_PORT="%(ENV_PLAYWRIGHT_MCP_PORT)s",CHROME_USER_DATA_DIR="%(ENV_CHROME_USER_DATA_DIR)s",DISPLAY="%(ENV_DISPLAY)s",REMOTE_DEBUG_PORT="9222",NODE_ENV="production"

